# crontab
*/30 * * * * alerta heartbeats --alert --severity major
